# addon/ya2dlna/services.d/dlna/run
#!/bin/sh
set -x
echo "$(date) Starting DLNA service" >> /tmp/debug.log
cd /app
exec python3 -m src.dlna_stream_server.main 2>&1