# addon/ya2dlna/rootfs/etc/services.d/dlna/run
#!/bin/sh
exec python3 -m src.dlna_stream_server.main